# etc/profile.d/appscale-tools.sh
export APPSCALE_TOOLS_PATH=/usr/local/appscale-tools
export PATH=${PATH}:${APPSCALE_TOOLS_PATH}/bin
